# jobs/crontab
0 0 * * * root python /app/jobs/daily_summary.py >> /var/log/cron.log 2>&1
